# systemd/mount-hdd.service
[Unit]
Description=Auto-mount external HDD for Picapport
After=network.target

[Service]
Type=oneshot
ExecStart=/home/pi/pi-photo-hub/scripts/mount-hdd.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
