# systemd/photo-api.service
[Unit]
Description=Photo API Server
After=network.target

[Service]
WorkingDirectory=/home/pi/pi-photo-hub/api
ExecStart=/usr/bin/node /home/pi/pi-photo-hub/api/server.js
Restart=always
User=pi
Environment=PATH=/opt/node/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
