# systemd/photo-api.service.template
[Unit]
Description=Photo API Server
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/pi-photo-hub/api
ExecStart=/opt/node/bin/node /home/pi/pi-photo-hub/api/server.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
