# systemd/picapport.service.template
[Unit]
Description=Picapport Photo Server
After=network.target

[Service]
User=pi
Group=pi
WorkingDirectory=/home/pi/pi-photo-hub/picapport
ExecStart=/opt/java/bin/java -Dpicapport.home=/home/pi/.picapport -Dpicapport.gui.enabled=false -jar picapport.jar
Restart=on-failure
RestartSec=5
Environment=DISPLAY=:0

[Install]
WantedBy=multi-user.target
