# initd/picapport
#!/bin/sh
### BEGIN INIT INFO
# Provides:          picapport
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:
# Default-Stop:
# Short-Description: Start Picapport via init.d (controlled by systemd wrapper)
# Description:       Headless Picapport daemon
### END INIT INFO

NAME=picapport
DESC="PicApport Daemon"
DAEMON="/opt/picapport/start-picapport.sh"
PIDFILE="/var/run/$NAME.pid"
USER="pi"

case "$1" in
  start)
    echo "Starting $DESC..."
    if [ -f "$PIDFILE" ] && kill -0 "$(cat "$PIDFILE")" 2>/dev/null; then
      echo "$NAME already running (PID $(cat "$PIDFILE"))"
      exit 0
    fi
    start-stop-daemon --start --background --make-pidfile --pidfile "$PIDFILE" --chuid "$USER" --exec "$DAEMON"
    sleep 1
    if [ -f "$PIDFILE" ] && kill -0 "$(cat "$PIDFILE")" 2>/dev/null; then
      echo "$NAME started (PID $(cat "$PIDFILE"))"
    else
      echo "Failed to start $NAME"
      exit 1
    fi
    ;;
  stop)
    echo "Stopping $DESC..."
    if [ -f "$PIDFILE" ]; then
      start-stop-daemon --stop --pidfile "$PIDFILE" --retry=TERM/30/KILL/5
      rm -f "$PIDFILE"
      echo "$NAME stopped."
    else
      echo "$NAME not running."
    fi
    ;;
  restart)
    $0 stop
    sleep 1
    $0 start
    ;;
  status)
    if [ -f "$PIDFILE" ] && kill -0 "$(cat "$PIDFILE")" 2>/dev/null; then
      echo "$NAME is running (PID $(cat "$PIDFILE"))"
      exit 0
    else
      echo "$NAME is not running"
      exit 1
    fi
    ;;
  *)
    echo "Usage: /etc/init.d/$NAME {start|stop|restart|status}"
    exit 2
    ;;
esac

exit 0
